# supervisor/queue-worker.ini
[program:messenger-consume]
command=php /path/to/your/app/bin/console messenger:consume async --time-limit=3600 ---memory-limit=128M --limit=5
user=juyal
numprocs=5
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d
